# deploy/growcontrol-webapi.service
[Unit]
Description=Growcontrol onboarding web API
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Growcontrol
EnvironmentFile=-/home/pi/Growcontrol/.env
ExecStart=/usr/bin/python3 /home/pi/Growcontrol/growcontrol_webapi.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
